# infra/kso-linux/systemd/kso-state-adapter.service
[Unit]
Description=KSO State Adapter — writes kso_state.json from UKM 4 terminal
Documentation=https://hermes-agent.nousresearch.com/docs

[Service]
Type=simple
EnvironmentFile=/etc/verny/kso/state-adapter.env
ExecStart=/usr/bin/python3 -m kso_state_adapter.cli daemon \
    --root /var/lib/verny/kso \
    --source-state ${VERNY_KSO_STATIC_STATE} \
    --health-file /run/verny/kso/state-adapter-health.json

WorkingDirectory=/opt/verny/kso
Restart=always
RestartSec=5

# ── Security hardening ────────────────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/verny/kso /run/verny/kso /var/log/verny/kso
ReadOnlyPaths=/opt/verny/kso /etc/verny/kso

# ── Logging ───────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kso-state-adapter

[Install]
WantedBy=multi-user.target
